# install-cyclescape-issuemap/cyclescape-issuemap.cron
# Update site each hour from Git
0 * * * * cyclestreets cd /var/www/cyclescape-issuemap/ && git pull --quiet
