# install-telluswhere/telluswhere.cron
# Update site each hour from Git
0 * * * * cyclestreets cd /var/www/telluswhere/ && git pull --quiet
