# live-deployment/cyclestreets-update.cron
# Update site regularly from Git, and run any database migrations
0,15,30,45 *  * * *  cyclestreets  cd /websites/www/content/ && git pull --quiet && yarn install --silent && composer install -n -q
0,15,30,45 *  * * *  cyclestreets  php /opt/cyclestreets-setup/live-deployment/cyclestreets-update-db.php
