# live-deployment/cyclestreets-setup-update.cron
# Update regularly from Git
0,15,30,45 *  * * *  cyclestreets  cd /opt/cyclestreets-setup && git pull --quiet
