# install-mobileweb/mobileweb.cron
# Update site regularly from Git
0,15,30,45 *  * * *  cyclestreets  cd /var/www/mobileweb/ && git pull --quiet
